# webservice/build/worker/services.d/messenger-consume-events/run
#!/command/with-contenv sh
php /app/bin/console messenger:consume async-event --time-limit=3600 --memory-limit=256M
