# webservice/build/worker/services.d/messenger-consume-cmds/run
#!/command/with-contenv sh
php /app/bin/console messenger:consume async-cmd --time-limit=3600 --memory-limit=256M
