# webservice/build/worker-intelpage/services.d/intelpage-daemon/run
#!/command/with-contenv sh
php /app/bin/console app:intel-page:run-daemon default-pager
